# var/lib/dpkg/info/certbot.list
/.
/etc
/etc/cron.d
/etc/cron.d/certbot
/etc/letsencrypt
/etc/letsencrypt/cli.ini
/etc/logrotate.d
/etc/logrotate.d/certbot
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/certbot.service
/lib/systemd/system/certbot.timer
/usr
/usr/bin
/usr/bin/certbot
/usr/share
/usr/share/doc
/usr/share/doc/certbot
/usr/share/doc/certbot/README.rst.gz
/usr/share/doc/certbot/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/certbot.1.gz
/usr/bin/letsencrypt
/usr/share/doc/certbot/changelog.Debian.gz
/usr/share/man/man1/letsencrypt.1.gz
