# var/lib/dpkg/info/certbot.conffiles
/etc/cron.d/certbot
/etc/letsencrypt/cli.ini
/etc/logrotate.d/certbot
