# var/lib/dpkg/info/ebtables.conffiles
/etc/default/ebtables
/etc/ethertypes
/etc/init.d/ebtables
