# var/lib/dpkg/info/ec2-hibinit-agent.list
/.
/etc
/etc/acpi
/etc/acpi/actions
/etc/acpi/actions/hibinit-power.sh
/etc/acpi/actions/sleep.sh
/etc/acpi/events
/etc/acpi/events/hibinit-power
/etc/acpi/events/hibinit-sleep
/etc/hibinit-config.cfg
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/hibinit-agent.service
/lib/systemd/system-sleep
/lib/systemd/system-sleep/hibinit-agent
/usr
/usr/bin
/usr/bin/hibinit-agent
/usr/lib
/usr/lib/ec2-hibinit-agent
/usr/lib/ec2-hibinit-agent/hibinit-resume
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/ec2_hibinit_agent-1.0.0.egg-info
/usr/lib/python3/dist-packages/ec2_hibinit_agent-1.0.0.egg-info/PKG-INFO
/usr/lib/python3/dist-packages/ec2_hibinit_agent-1.0.0.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/ec2_hibinit_agent-1.0.0.egg-info/top_level.txt
/usr/lib/systemd
/usr/lib/systemd/logind.conf.d
/usr/lib/systemd/logind.conf.d/ec2-hibinit-agent-ignore-powerkey.conf
/usr/share
/usr/share/doc
/usr/share/doc/ec2-hibinit-agent
/usr/share/doc/ec2-hibinit-agent/README
/usr/share/doc/ec2-hibinit-agent/TODO.Debian
/usr/share/doc/ec2-hibinit-agent/changelog.Debian.gz
/usr/share/doc/ec2-hibinit-agent/copyright
