# var/lib/dpkg/info/hibagent.prerm
#!/bin/sh
set -e



# Manually added instead of using dh_installinit
if [ -x "/etc/init.d/hibagent" ]; then
        invoke-rc.d hibagent stop || true
fi
