# var/lib/dpkg/info/hibagent.conffiles
/etc/hibagent-config.cfg
/etc/init.d/hibagent
