# var/lib/dpkg/info/open-vm-tools.list
/.
/etc
/etc/init.d
/etc/init.d/open-vm-tools
/etc/pam.d
/etc/pam.d/vmtoolsd
/etc/vmware-tools
/etc/vmware-tools/poweroff-vm-default
/etc/vmware-tools/poweron-vm-default
/etc/vmware-tools/resume-vm-default
/etc/vmware-tools/scripts
/etc/vmware-tools/scripts/vmware
/etc/vmware-tools/scripts/vmware/network
/etc/vmware-tools/statechange.subr
/etc/vmware-tools/suspend-vm-default
/etc/vmware-tools/tools.conf
/etc/vmware-tools/tools.conf.example
/etc/vmware-tools/vgauth
/etc/vmware-tools/vgauth/schemas
/etc/vmware-tools/vgauth/schemas/XMLSchema-hasFacetAndProperty.xsd
/etc/vmware-tools/vgauth/schemas/XMLSchema-instance.xsd
/etc/vmware-tools/vgauth/schemas/XMLSchema.dtd
/etc/vmware-tools/vgauth/schemas/XMLSchema.xsd
/etc/vmware-tools/vgauth/schemas/catalog.xml
/etc/vmware-tools/vgauth/schemas/datatypes.dtd
/etc/vmware-tools/vgauth/schemas/saml-schema-assertion-2.0.xsd
/etc/vmware-tools/vgauth/schemas/xenc-schema.xsd
/etc/vmware-tools/vgauth/schemas/xml.xsd
/etc/vmware-tools/vgauth/schemas/xmldsig-core-schema.xsd
/etc/vmware-tools/vgauth.conf
/etc/vmware-tools/vm-support
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/open-vm-tools.service
/lib/systemd/system/vgauth.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-open-vm-tools.rules
/lib/udev/rules.d/99-vmware-scsi-udev.rules
/sbin
/sbin/mount.vmhgfs
/usr
/usr/bin
/usr/bin/VGAuthService
/usr/bin/vmhgfs-fuse
/usr/bin/vmtoolsd
/usr/bin/vmware-checkvm
/usr/bin/vmware-hgfsclient
/usr/bin/vmware-namespace-cmd
/usr/bin/vmware-rpctool
/usr/bin/vmware-toolbox-cmd
/usr/bin/vmware-vgauth-cmd
/usr/bin/vmware-vgauth-smoketest
/usr/bin/vmware-vmblock-fuse
/usr/bin/vmware-xferlogs
/usr/lib
/usr/lib/libDeployPkg.so.0.0.0
/usr/lib/libguestlib.so.0.0.0
/usr/lib/libhgfs.so.0.0.0
/usr/lib/libvgauth.so.0.0.0
/usr/lib/libvmtools.so.0.0.0
/usr/lib/open-vm-tools
/usr/lib/open-vm-tools/plugins
/usr/lib/open-vm-tools/plugins/common
/usr/lib/open-vm-tools/plugins/common/libhgfsServer.so
/usr/lib/open-vm-tools/plugins/common/libvix.so
/usr/lib/open-vm-tools/plugins/vmsvc
/usr/lib/open-vm-tools/plugins/vmsvc/libdeployPkgPlugin.so
/usr/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
/usr/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
/usr/lib/open-vm-tools/plugins/vmsvc/libresolutionKMS.so
/usr/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
/usr/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
/usr/sbin
/usr/share
/usr/share/doc
/usr/share/doc/open-vm-tools
/usr/share/doc/open-vm-tools/changelog.Debian.gz
/usr/share/doc/open-vm-tools/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/open-vm-tools
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/vmware-checkvm.1.gz
/usr/share/man/man1/vmware-toolbox-cmd.1.gz
/usr/share/man/man1/vmware-xferlogs.1.gz
/usr/share/man/man3
/usr/share/man/man3/libguestlib.3.gz
/usr/share/man/man3/libvmtools.3.gz
/usr/share/man/man9
/usr/share/man/man9/vmxnet.9.gz
/usr/share/open-vm-tools
/usr/share/open-vm-tools/messages
/usr/share/open-vm-tools/messages/de
/usr/share/open-vm-tools/messages/de/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/de/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/de/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/de/vmtoolsd.vmsg
/usr/share/open-vm-tools/messages/de/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/en
/usr/share/open-vm-tools/messages/en/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/en/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/en/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/es
/usr/share/open-vm-tools/messages/es/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/es/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/es/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/fr
/usr/share/open-vm-tools/messages/fr/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/fr/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/fr/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/it
/usr/share/open-vm-tools/messages/it/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/it/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/it/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/ja
/usr/share/open-vm-tools/messages/ja/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/ja/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/ja/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/ja/vmtoolsd.vmsg
/usr/share/open-vm-tools/messages/ja/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/ko
/usr/share/open-vm-tools/messages/ko/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/ko/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/ko/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/ko/vmtoolsd.vmsg
/usr/share/open-vm-tools/messages/ko/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/zh_CN
/usr/share/open-vm-tools/messages/zh_CN/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/zh_CN/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/zh_CN/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/zh_CN/vmware-vgauth-cmd.vmsg
/usr/share/open-vm-tools/messages/zh_TW
/usr/share/open-vm-tools/messages/zh_TW/VGAuthService.vmsg
/usr/share/open-vm-tools/messages/zh_TW/libvgauth.la.vmsg
/usr/share/open-vm-tools/messages/zh_TW/vmware-vgauth-cmd.vmsg
/usr/lib/libDeployPkg.so.0
/usr/lib/libguestlib.so.0
/usr/lib/libhgfs.so.0
/usr/lib/libvgauth.so.0
/usr/lib/libvmtools.so.0
